# src/mumble/build/modules/mumble/files/etc/supervisor/conf.d/mumble.conf
[program:mumble]
command = /usr/sbin/murmurd -ini /etc/mumble-server.ini
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
