# src/mumble/build/modules/mumble/files/etc/supervisor/conf.d/mumble.conf
[program:mumble]
command = /usr/sbin/murmurd -fg -ini /etc/mumble-server.ini
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
