# src/mumble/run/modules/mumble/templates/mumble.conf.erb
[program:mumble]
command = /usr/sbin/murmurd -ini /etc/mumble-server.ini -supw <%= @supw %>
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
